# teetime-daily.service
[Unit]
Description=TeetimeFinder - Send daily status digest
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jon
WorkingDirectory=/srv/apps/teetime_scraper
ExecStart=/srv/apps/teetime_scraper/.venv/bin/python /srv/apps/teetime_scraper/scraper.py --daily-digest
StandardOutput=journal
StandardError=journal
